# deployment/llm-audit-cleanup.timer
[Unit]
Description=Run LLM Audit cleanup daily
Requires=llm-audit-cleanup.service

[Timer]
OnCalendar=daily
OnCalendar=03:00
Persistent=true

[Install]
WantedBy=timers.target


